# fuel_logger/rootfs/etc/services.d/fuel-logger/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Fuel Logger service"

export OPENAI_API_KEY="$(bashio::config 'OPENAI_API_KEY')"
export GOOGLE_CLIENT_EMAIL="$(bashio::config 'GOOGLE_CLIENT_EMAIL')"
export GOOGLE_PRIVATE_KEY="$(bashio::config 'GOOGLE_PRIVATE_KEY')"
export GOOGLE_SHEET_ID="$(bashio::config 'GOOGLE_SHEET_ID')"

# Respect the configured port from the add-on options
PORT="$(bashio::addon.port 3000)"
export PORT
bashio::log.info "Listening on port ${PORT}"

exec npm --prefix /app/backend start
